# supervisor/conf.d/tailscaled.conf
[program:tailscaled]
command=/usr/sbin/tailscaled --state=/workspace/.piclaw/tailscale/tailscaled.state --socket=/workspace/.piclaw/tailscale/tailscaled.sock
user=root
autostart=false
autorestart=true
startsecs=3
startretries=10
stopsignal=INT
stopwaitsecs=10
stdout_logfile=/var/log/tailscale/tailscaled.stdout.log
stderr_logfile=/var/log/tailscale/tailscaled.stderr.log
